# supervisor/dockyard.conf
[program:run-nginx]
directory=/opt/dockyard/nginx
command=python run_nginx.py
autorestart=false
startsecs=0

[program:run-registry]
directory=/opt/dockyard/registry
command=python run_registry.py
